# daemon/googolplex-theater.conf
description "googolplex-theater"
author  "yuzawa-san"

start on runlevel [2345]
stop on shutdown

respawn
respawn limit 5 20

exec /usr/bin/java -jar /path/to/you/googolplex-theater-all.jar --cast-config /path/to/your/config.json