# src/dist/service/com.jyuzawa.googolplex-theater.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC -//Apple Computer//DTD PLIST 1.0//EN
           http://www.apple.com/DTDs/PropertyList-1.0.dtd >
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.jyuzawa.googolplex-theater</string>
    <key>ProgramArguments</key>
    <array>
      <string>/path/to/your/googolplex-theater-VERSION/bin/googolplex-theater</string>
    </array>
    <key>KeepAlive</key>
    <true />
    <key>StandardOutPath</key>
    <string>/Library/Logs/googolplex-theater/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Library/Logs/googolplex-theater/stderr.log</string>
  </dict>
</plist>