# src/dist/service/googolplex-theater.conf
description "googolplex-theater"
author  "yuzawa-san"

start on runlevel [2345]
stop on shutdown

respawn
respawn limit 5 20

exec /path/to/your/googolplex-theater-VERSION/bin/googolplex-theater