# src/dist/service/googolplex-theater.service
[Unit]
Description=googolplex-theater
After=network.target

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
User=root
Restart=on-failure
RestartSec=1
StartLimitInterval=60
StartLimitBurst=5
ExecStart=/path/to/your/googolplex-theater-VERSION/bin/googolplex-theater
