# deb/postrm
#!/bin/sh

debsystemctl=$(command -v deb-systemd-invoke || echo systemctl)
$debsystemctl disable googolplex-theater
